# .config/profile.d/aliases.sh
###############################################################################
# ls                                                                          #
###############################################################################
alias ls='ls --color=tty'

###############################################################################
# Git                                                                         #
###############################################################################
alias gs='git status'
alias ga='git add'
alias gp='git push'
alias gc='git commit'
alias gl='git log'
alias gpl='git pull'

alias git-dot='/usr/bin/git --git-dir=$HOME/.config/dotfiles --work-tree=$HOME'

###############################################################################
# Nvidia                                                                      #
###############################################################################
alias nsmi='nvidia-smi'
alias nsmi-u='nvidia-htop.py'

###############################################################################
# Tar                                                                         #
###############################################################################
alias tgz='tar -czvf'
alias untgz='tar -xvzf'

###############################################################################
# Disk Usage                                                                  #
###############################################################################
alias dua='du -hs $(ls -A) | sort -h'
