# .config/profile.d/pixi.sh
export PIXI_HOME="$HOME/.local/share/pixi"
export PATH="$PIXI_HOME/bin:$PATH"
