# .config/profile.d/vim.sh
export VIMINIT="set nocp | source ${XDG_CONFIG_HOME}/vim/vimrc"
export VIMDOTDIR="$XDG_CACHE_HOME/vim"
